# ts_ws/devel/etc/catkin/profile.d/villa_audio.sh
/home/tracy/catkin_ws/devel/.private/villa_audio/etc/catkin/profile.d/villa_audio.sh